# graphite/carbon-cache.conf
[program:carbon-cache]
numprocs=1
autostart=true
startretries=2
autorestart=unexpected
user=root
redirect_stderr=true
directory=/opt/graphite/bin/
stdout_logfile=/var/log/supervisor/%(program_name)s-stdout.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/%(program_name)s-stderr.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
stderr_capture_maxbytes=1MB
command=python carbon-cache.py --debug start
stopsignal=QUIT